# tecnotime.service
[Unit]
Description=TecnoTime - Sistema de Checadores ZKTeco
After=network.target mysql.service

[Service]
Type=simple
User=ccomputo
Group=ccomputo
WorkingDirectory=/home/ccomputo/projects/rino
Environment="PATH=/home/ccomputo/projects/rino/.venv/bin"
ExecStart=/home/ccomputo/projects/rino/.venv/bin/gunicorn --bind 0.0.0.0:5000 --workers 4 --timeout 120 --access-logfile /var/log/tecnotime/access.log --error-logfile /var/log/tecnotime/error.log main:app

# Reiniciar automáticamente si falla
Restart=always
RestartSec=10

# Limitar recursos
LimitNOFILE=4096

# Variables de entorno
EnvironmentFile=/home/ccomputo/projects/rino/.env

[Install]
WantedBy=multi-user.target
